# com.hermes.router-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.router-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update YOUR_USERNAME to match your macOS username -->
        <string>/Users/YOUR_USERNAME/.hermes/hermes-agent/venv/bin/python3</string>
        <string>/Users/YOUR_USERNAME/.hermes/skills/router_server.py</string>
        <string>--port</string>
        <string>5050</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/.hermes</string>

    <!-- Start automatically when user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if the process crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.hermes/logs/router_proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.hermes/logs/router_proxy.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.hermes/hermes-agent/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Prevent rapid restart loops on repeated crashes -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
